# ttrss/ttrss/root/etc/s6/ttrss-update-daemon/run
#!/bin/sh

USER=nobody
exec su-exec $USER php /var/www/ttrss/update_daemon2.php
